# knxd/rootfs/etc/s6-overlay/s6-rc.d/init-knxd-config/run
#!/usr/bin/with-contenv bashio

readonly CONF="/config/knxd.ini"
readonly CONF_TEMPLATE="/etc/knxd.ini"
readonly TOKEN="#@autogenerated"

# =======================================
# custom configuration file
# =======================================
if bashio::fs.file_exists "${CONF}"; then
	# check if existing file starts with our token (indicating auto-update) or if it is a custom file
	read -r FILEHEADER < "${CONF}"

	if ! [ "${TOKEN}" = "${FILEHEADER}" ]; then
		bashio::log.warning "Found custom configuration file!"
		bashio::log.info "Any configuration from the add-ons configuration section will be ignored."
		bashio::log.info "Using custom configuration:"
		bashio::log.info "$(<$CONF)"
		exit 0
	fi

fi


# =======================================
# custom configuration setting
# =======================================

# first check if a custom config is specified and use that
if bashio::config.has_value 'custom_config' && [ $(bashio::string.length "$(bashio::config 'custom_config')") > 5 ]; then
  bashio::log.warning "Deprecated config option 'custom_config' is used. This option will be removed in a future update. Your configuration is automatically migrated and moved to the knxd.ini file in the addon_config folder of this add-on."
	bashio::log.info "Using custom configuration:"
	bashio::log.info "$(bashio::config 'custom_config')"
	cat > ${CONF} << EOF
$(bashio::config 'custom_config')
EOF
	exit 0
fi

# =======================================
# auto configuration
# =======================================

# check if an interface is defined
if ! $(bashio::config.has_value 'interface'); then
	bashio::exit.nok "An interface type needs to be configured"
fi

# check if a device is selected, which is mandatory for anything but USB
if $(bashio::config.equals 'interface' 'usb') && ! $(bashio::config.has_value 'device'); then
	bashio::exit.nok "A device needs to be configured"
fi

# in case of USB interface, we can not use the device path directly
# because knxd expects us to pass the bus ID and device ID separately
# so this is what we do now
declare USB_DEVICE='';
declare USB_BUS='';
declare device=$(bashio::config 'device')

if $(bashio::config.equals 'interface' 'usb') && [ $(bashio::string.substring "${device}" 0 13) == "/dev/bus/usb/" ]; then
	USB_BUS=$(echo "${device}" | cut -d / -f 5 | sed 's/^0*//g')
	USB_DEVICE=$(echo "${device}" | cut -d / -f 6 | sed 's/^0*//g')
fi

# create config file from template
cp "${CONF_TEMPLATE}" "${CONF}"

# replace placeholders and write it back
sed -i "s|%%TOKEN%%|${TOKEN}|g" "${CONF}"
sed -i "s|%%INTERFACE%%|$(bashio::config 'interface')|g" "${CONF}"
sed -i "s|%%DEVICE%%|$(bashio::config 'device')|g" "${CONF}"
sed -i "s|%%ADDRESS%%|$(bashio::config 'address')|g" "${CONF}"
sed -i "s|%%CLIENT_ADDRESS%%|$(bashio::config 'client_address')|g" "${CONF}"
sed -i "s|%%FILTERS%%|$(bashio::config 'usb_filters')|g" "${CONF}"
sed -i "s|%%DEBUG_ERROR_LEVEL%%|$(bashio::config 'log_error_level')|g" "${CONF}"
sed -i "s|%%USB_DEVICE%%|${USB_DEVICE}|g" "${CONF}"
sed -i "s|%%USB_BUS%%|${USB_BUS}|g" "${CONF}"
sed -i "s|%%IP_ADDRESS%%|$(bashio::config 'ip_address')|g" "${CONF}"
sed -i "s|%%DEST_PORT%%|$(bashio::config 'dest_port')|g" "${CONF}"


# =======================================
# log output configuration
# =======================================
bashio::log.info "Interface:      $(bashio::config 'interface')"
bashio::log.info "Device:         $(bashio::config 'device')"
bashio::log.info "Address:        $(bashio::config 'address')"
bashio::log.info "Client Address: $(bashio::config 'client_address')"
bashio::log.info "Filters:        $(bashio::config 'usb_filters')"
bashio::log.info "Log-Level:      $(bashio::config 'log_error_level')"
if [ "${USB_DEVICE}" != "" ]; then
  bashio::log.info "USB-Device ID:  ${USB_DEVICE}"
  bashio::log.info "USB-Bus ID:     ${USB_BUS}"
fi
bashio::log.info "IP-Address:     $(bashio::config 'ip_address')"
bashio::log.info "Dest-Port:      $(bashio::config 'dest_port')"
