# knxd/rootfs/etc/services.d/knxd/run
#!/usr/bin/with-contenv bash

exec knxd /etc/knxd.ini