# knxd/rootfs/etc/s6-overlay/s6-rc.d/svc-knxd/run
#!/usr/bin/with-contenv bash

exec knxd /config/knxd.ini
